# profile.d/40-rust.sh
# if cargo is installed, then add its environment to the path
if [ -f ~/.cargo/env ] ; then
    . ~/.cargo/env
fi

