# profile.d/77-music.sh
export MPD_HOST=/run/user/1000/mopidy
alias tj='mpc add https://uncertain.fm/streams/sonos && mpc play'
