# profile.d/85-no-home-cruft.sh
# remove / relocate excessive nonsense from ~

# less
rm -f ~/.lesshst
export LESSHISTFILE=$HOME/.cache/lesshst

# wget
export WGETRC=$HOME/.config/wgetrc
