# profile.d/90-gpg-agent.sh
# always favor the gpg-agent's ssh socket
export SSH_AUTH_SOCK="/run/user/$(id -u)/gnupg/S.gpg-agent.ssh"
