# config/systemd/user/gpg-agent.service.d/override.conf
[Service]
Environment=GNUPGHOME=/home/benferse/.config/gnupg
